# daemon/units-daemon.service
[Unit]
Description=Units Server

[Service]
ExecStart=/usr/sbin/units-daemon

[Install]
WantedBy=multi-user.target

